# systemd/myplant-bot.service
[Unit]
Description=myPlant FastAPI Telegram Bot
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/myPlant_bot
EnvironmentFile=/home/ec2-user/myPlant_bot/.env
Environment=APP_DIR=/home/ec2-user/myPlant_bot
ExecStart=/bin/bash /home/ec2-user/myPlant_bot/scripts/start_server.sh
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
